# docker/service/hhvm/run
#!/bin/bash
exec /usr/bin/hhvm --mode server --user www-data --config /etc/hhvm/server.ini
